# examples/db-backup.timer
[Unit]
Description=Daily Backup of MySQL data

[Timer]
OnBootSec=60min
OnUnitActiveSec=1d
